# yakk/templates/systemd/yakk-serve.service
# yakk-serve.service v1.0.0
# Last updated: 2026-01-21
# Simplified: start script handles config via yakk.env

[Unit]
Description=Yakk HTTP Server (Streamable HTTP/SSE)
After=network.target

[Service]
Type=simple
ExecStart={START_SCRIPT}
Restart=on-failure
RestartSec=10

# Environment - start script sources yakk.env for host/port/transport config
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yakk-serve

[Install]
WantedBy=default.target
